# dobrinya-bot.service
[Unit]
Description=Dobrinya VPN Telegram Bot v14
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/dobrinya_bot
ExecStart=/root/dobrinya_bot/venv/bin/python bot.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
